# modules/core/templates/rc.local.erb
#!/bin/bash
# /etc/rc.local - run once at boot time
# Put any local setup commands in here:
